# cookbook/templates/default/upstart/met-api.conf.erb
description     "Meteorological API"
author          "MET Norway <post@met.no>"

start on runlevel [2345]
exec <%= node['met-api']['server_binary'] %> \
    -Dhttp.port=<%= node['met-api']['listen_port'] %> \
    -Dhttp.address=<%= node['met-api']['listen_address'] %> \
    -Djava.net.preferIPv4Stack=true \
    -Dconfig.file=/etc/met-api/application.conf
respawn
